# hadoop/hbase/images/supervisord.conf
[supervisord]
nodaemon=true

[program:start-k8s-hbase]
command=/hbase/bin/start-k8s-hbase.sh
