# hadoop/hdfs/namenode/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop-hdfs-namenode]
command=/sbin/hadoop-daemon.sh start namenode


