# hadoop/hdfs/datanode/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop-hdfs-datanode]
command=/sbin/hadoop-daemon.sh start datanode
